# scenarios/sigmahq/proc_creation_win_schtasks_appdata_local_system.yml
type: trigger
name: sigmahq/proc_creation_win_schtasks_appdata_local_system
description: |
  Detects the creation of a schtask that executes a file from C:\Users\<USER>\AppData\Local
filter: |
  (evt.Meta.service == 'sysmon' && evt.Parsed.EventID == '1') && (evt.Parsed.Image endsWith '\\schtasks.exe' && evt.Parsed.CommandLine contains '/Create' && evt.Parsed.CommandLine contains '/RU' && evt.Parsed.CommandLine contains '/TR' && evt.Parsed.CommandLine contains 'C:\\Users\\' && evt.Parsed.CommandLine contains '\\AppData\\Local\\' && (evt.Parsed.CommandLine contains 'NT AUT' || evt.Parsed.CommandLine contains ' SYSTEM ') && not (evt.Parsed.ParentImage contains '\\AppData\\Local\\Temp\\' && evt.Parsed.ParentImage contains 'TeamViewer_.exe' && evt.Parsed.Image endsWith '\\schtasks.exe' && evt.Parsed.CommandLine contains '/TN TVInstallRestore'))
blackhole: 2m
#status: test
labels:
  service: windows
  confidence: 1
  spoofable: 0
  classification:
   - attack.t1053.005
   - attack.t1059.001

  label: "Suspicious Schtasks Execution AppData Folder"
  behavior : "windows:audit"
  remediation: false

scope:
  type: ParentProcessId
  expression: evt.Parsed.ParentProcessId

